# blink_ip.timer
[Unit]
Description=runs blink_ip service every ten seconds

[Timer]
OnBootSec=60sec
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target
