# provision.timer
[Unit]
Description=kicks provision service after boot
Wants=network-online.target
After=network-online.target

[Timer]
OnBootSec=60

[Install]
WantedBy=timers.target
