# provision.timer
[Unit]
Description=kicks provision service after boot

[Timer]
OnBootSec=10

[Install]
WantedBy=timers.target
