# provision.timer
[Unit]
Description=runs provision service weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
